# overlay/etc/cron.hourly/osqa-email-alerts
#!/bin/bash -e
# send osqa email alerts

su www-data -c "/var/www/osqa/manage.py send_email_alerts"

